# system/skeleton/etc/inittab
::sysinit:/etc/init.d/rc sysinit
::shutdown:/etc/init.d/rc shutdown

::restart:/sbin/init
# Put a getty on the serial port
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100 # GENERIC_SERIAL
