# system/skeleton/etc/inittab
::sysinit:/etc/init.d/rc sysinit
::shutdown:/etc/init.d/rc shutdown

::restart:/sbin/init
::respawn:/sbin/getty 115200 tty1
