# package/logrotate/logrotate.cron
#!/bin/bash
logrotate /etc/logrotate.conf &> /dev/null
