# package/fake-hwclock/fake-hwclock.cron
#!/bin/bash
fake-hwclock save &> /dev/null
