# package/busybox/ntpd.cron
# NTPD keepalive cron task
*/5 * * * * /etc/init.d/ntpd keepalive
